# _tools/launchagents/com.tomoko.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tomoko.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/seijiro/Sync/sync_work/by-llms/tomoko/_tools/run_daily_launchagent.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/seijiro/Sync/sync_work/by-llms/tomoko</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/tomoko-daily.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/tomoko-daily.launchd.err.log</string>
</dict>
</plist>
